# supervisord.conf
[supervisord]
nodaemon=true
user=nginx
logfile=/var/log/supervisor/supervisord.log ; (main log file; default $CWD/supervisord.log)

[program:php-fpm]
command=php-fpm83 -F
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/php-fpm.err.log
stdout_logfile=/var/log/supervisor/php-fpm.out.log

[program:nginx]
command=nginx -g 'daemon off;'
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log

[program:fcgiwrap]
command=spawn-fcgi -n -s /var/run/fcgiwrap.socket -P /var/run/fcgiwrap.pid /usr/bin/fcgiwrap
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/fcgiwrap.err.log
stdout_logfile=/var/log/supervisor/fcgiwrap.out.log

[program:nagios]
command=/usr/local/nagios/bin/nagios /usr/local/nagios/etc/nagios.cfg
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/nagios.err.log
stdout_logfile=/var/log/supervisor/nagios.out.log